# deployment/crypto-trader-dashboard.service
[Unit]
Description=Crypto AI Trader Dashboard (Streamlit)
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=jetson
WorkingDirectory=/home/jetson/crypto_ai_trader
ExecStart=/home/jetson/crypto_ai_trader/.venv/bin/python -m streamlit run /home/jetson/crypto_ai_trader/dashboard/streamlit_app.py --server.address 0.0.0.0 --server.port 8501 --browser.gatherUsageStats false
EnvironmentFile=/home/jetson/crypto_ai_trader/.env
Environment=PYTHONUNBUFFERED=1
Restart=on-failure
RestartSec=15s
KillSignal=SIGINT
TimeoutStopSec=30s

# Keep the dashboard lighter than the runtime worker on Nano
MemoryMax=500M
MemorySwapMax=1G

# Logging via journald (view with: journalctl -fu crypto-trader-dashboard)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=crypto-trader-dashboard

[Install]
WantedBy=multi-user.target
